# arch-specifics/dot_config/autostart/nitrogen.desktop
[Desktop Entry]
Exec=nitrogen --restore
Icon=xscreensaver
Terminal=false
Name=Screensaver
Type=Application
