# arch-specifics/dot_config/autostart/nmapplet.desktop
[Desktop Entry]
Exec=nm-applet
Icon=nm-applet
Terminal=false
Name=NetworkManagerApplet
Type=Application
