# arch-specifics/dot_config/autostart/redshift.desktop
[Desktop Entry]
Exec=redshift -t 6000:1300 -l -33.95:151.04
Icon=redshift
Terminal=false
Name=Redshift
Type=Application
