# arch-specifics/dot_config/autostart/redshift.desktop
[Desktop Entry]
Exec=redshift -t 6500k:2700k
Icon=redshift
Terminal=false
Name=Redshift
Type=Application
