# arch-specifics/dot_config/autostart/twmnd.desktop
[Desktop Entry]
Exec=twmnd
Icon=twmnd
Terminal=false
Name=Notifications
Type=Application
