# arch-specifics/dot_config/autostart/udiskie.desktop
[Desktop Entry]
Exec=udiskie --no-automount --tray --use-udisks2
Icon=udiskie
Terminal=false
Name=Udiskie
Type=Application
