# arch-specifics/dot_config/autostart/xscreensaver.desktop
# [Desktop Entry]
# Exec=xscreensaver -nosplash
# Icon=xscreensaver
# Terminal=false
# Name=Screensaver
# Type=Application
